# docker/supervisor-smtp-server.conf
[program:smtp-server]
command=/opt/venv/bin/python manage.py run_smtp_server
directory=/app
autostart=false
autorestart=true
stdout_logfile=/app/logs/smtp-server-stdout.log
stderr_logfile=/app/logs/smtp-server-stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
user=spwig
priority=15
startsecs=5
stopwaitsecs=10
environment=PYTHONUNBUFFERED="1",PATH="/opt/venv/bin:%(ENV_PATH)s"

# Note: autostart=false because SMTP server should only start when
# a built-in SMTP email account is configured via the admin wizard.
# Admins can start it via: supervisorctl start smtp-server
